# maquinaexterna/peticionmaquinaexterna.cron
* * * * * /usr/share/easy2train/peticionmaquinaexterna.sh

